# ansible/roles/nextcloud/templates/nextcloud.cron.j2

*/5 * * * * www-data php -f {{ nextcloud_config.install_dir }}/nextcloud/cron.php

